# bookworm.service
[Unit]
Description=Bookworm — Jean's Book Club
After=network.target

[Service]
Type=simple
User=rod
WorkingDirectory=/opt/bookworm
EnvironmentFile=/opt/bookworm/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bookworm

[Install]
WantedBy=multi-user.target
